# roles/service/redis/files/rc.local
#!/bin/bash
./etc/rc.d/init.d/redis start &
touch /var/lock/subsys/local